# scripts/daytradebot.service
[Unit]
Description=DayTradeBot intraday loop (Alpaca paper/live)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/stonk/projects/DayTradeBot
EnvironmentFile=-/home/stonk/projects/DayTradeBot/.env
Environment=DAYTRADEBOT_LOG_FILE=/home/stonk/projects/DayTradeBot/logs/daytradebot.log
ExecStart=/home/stonk/projects/DayTradeBot/scripts/run-daemon.sh
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=600
StartLimitBurst=10

[Install]
WantedBy=default.target
